# .ddev/web-build/autopublish.cron
*/1 * * * * cd /var/www/html && php site/plugins/kirby-autopublish/worker.php
